# debian/etc/dehydrated/cron/dehydrated.cron
#!/bin/bash

cd /etc/dehydrated
/usr/sbin/dehydrated --cron

