# debian/conf/dehydrated.cron
0 9 * * 6 /usr/local/bin/dehydrated -c -f /etc/letsencrypt/config
